# packaging/systemd/puppycam.service
[Service]
Type=simple
User=puppycam
Group=video

# Optional env file; service still runs with defaults below
EnvironmentFile=-/etc/puppycam.env

# Defaults (used if not set in /etc/puppycam.env)
Environment=PUPPYCAM_DEVICE=/dev/video0
Environment=PUPPYCAM_WIDTH=1280
Environment=PUPPYCAM_HEIGHT=720
Environment=PUPPYCAM_FPS=10
Environment=PUPPYCAM_PORT=8080

ExecStart=/usr/local/bin/puppycam \
  --device ${PUPPYCAM_DEVICE} \
  --width ${PUPPYCAM_WIDTH} \
  --height ${PUPPYCAM_HEIGHT} \
  --fps ${PUPPYCAM_FPS} \
  --port ${PUPPYCAM_PORT}

Restart=on-failure
RestartSec=2

NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full

StandardOutput=journal
StandardError=journal
