# packaging/systemd/puppycam.service
[Unit]
Description=PuppyCam (Qt6 V4L2 MJPEG streamer)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=puppycam
Group=video

# Run on localhost only (Tailscale Serve will expose it)
Environment=PUPPYCAM_DEVICE=/dev/v4l/by-id/usb-EMEET_EMEET_SmartCam_C960_A250926000113128-video-index0

ExecStart=/usr/local/bin/puppycam \
  --device ${PUPPYCAM_DEVICE} \
  --width 1280 --height 720 --fps 10 --port 8080

Restart=on-failure
RestartSec=2

# Basic hardening (safe defaults; we can tighten later)
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
